# xdg-config/autostart/redshift-gtk.desktop
[Desktop Entry]
Terminal=false
X-GNOME-Autostart-enabled=true
GenericName=Color temperature adjustment
Name[de]=Redshift
Type=Application
GenericName[da]=Justering af farvetemperatur
GenericName[cs]=Přizpůsobení barevné teploty
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
Categories=Utility;
Name[lt]=Redshift
Name[pt_BR]=Redshift
Comment[it]=Strumento per la regolazione della temperatura del colore
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
Name[pl]=Redshift
GenericName[it]=Regolazione della temperatura del colore
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Name[it]=Redshift
GenericName[tr]=Renk sıcaklığı ayarı
Icon=redshift
Name[tr]=Redshift
Name[hu]=Redshift
Name[ru]=Redshift
Comment[es]=Herramienta para el ajuste de la temperatura de color
GenericName[zh_CN]=色温调节
StartupNotify=true
Comment[de]=Farbtemperatur Anpassungswerkzeug
GenericName[pl]=Dostosowanie temperatury barwowej
Name[es]=Redshift
Name[zh_CN]=红移
Comment[zh_CN]=色温调节工具
GenericName[de]=Farbtemperatur Anpassung
Comment[tr]=Renk sıcaklığı ayarlama aracı
Hidden=false
Comment[fr]=Outil de réglage de la température de couleur
Name=Redshift
Version=1.0
GenericName[ru]=Настройка цветовой температуры
Comment[hu]=Színhőmérséklet beállító eszköz
GenericName[nl]=Kleurtemperatuur aanpassing
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
Name[fr]=Redshift
GenericName[lt]=Spalvos temperatūros reguliavimas
Exec=redshift-gtk
GenericName[fr]=Réglage de la température de couleur
GenericName[hu]=Szín hőmérséklet beállítás
Comment[da]=Justeringsværktøj for farvetemperatur
Comment=Color temperature adjustment tool
GenericName[es]=Ajuste de la temperatura de color
Comment[ru]=Инструмент регулирования цветовой температуры

