# cron/alpha.crontab
# crontab -- repeating background tasks.
MAILTO=g
PATH=/home/$LOGNAME/bin:/usr/bin:/usr/sbin:/usr/local/bin:/sbin:/bin
GKCRON=yes

# min	hour	d/month	month	d/week	command
*/5	*	*	*	*	/bin/sh $HOME/cf/lib/smtpq-notify.sh
*/5	*	*	*	*	/usr/bin/dconf dump / > $HOME/.config/dconf/user.dump
# end
