# cron/alpha.crontab
# crontab -- repeating background tasks. -*- mode: conf-space -*-
MAILTO=g
PATH=/home/$LOGNAME/bin:/usr/bin:/usr/sbin:/usr/local/bin:/sbin:/bin
MYLIB=$HOME/cf/lib
GKCRON=yes

# min	hour	d/month	month	d/week	command
*/5	*	*	*	*	/bin/sh $HOME/cf/lib/smtpq-notify.sh
*/1	*	*	*	*	/bin/sh $HOME/cf/lib/batmon.sh
# end
